# root/etc/services.d/rutorrent/run
#!/usr/bin/with-contenv bash
rtorrent s6-setuidgid abc /usr/bin/rtorrent -n -o import=/config/rtorrent/rtorrent.rc
